# consul-template/consul-template.ini
[program:consul-template]
command=consul-template -template "/usr/local/job/genrateServer.tpl:/usr/local/job/nginxServer.conf"  -consul-addr 10.23.27.87:80
directory=/usr/local/job
autorestart=true
autostart=true
stderr_logfile=/data/bigdata/consul-template/template.err.log
stdout_logfile=/data/bigdata/consul-template/template.out.log
user=root
stopsignal=INT
startsecs=1
