# consul-template/go.ini
[program:richive]
command=/usr/local/job/main -c /usr/local/job/openapi.yaml -e online
directory=/usr/local/job
autorestart=true
autostart=true
stderr_logfile=/data/bigdata/go/richive.err.log
stdout_logfile=/data/bigdata/go/richive.out.log
user=root
stopsignal=INT
startsecs=1
